# crontab
MAILTO="zpzhang@berkeley.edu"
# m h  dom mon dow   command
00 03 * * * /usr/bin/rsnapshot daily
00 01 * * 6 /usr/bin/rsnapshot weekly
00 00 1 * * /usr/bin/rsnapshot monthly

